# edgedbpkg/edgedb/IDENTIFIER-SLOT.service.in
[Unit]
Description=@@{description}
Documentation=@@{documentation}
After=syslog.target
After=network.target

[Service]
Type=notify

User=@@{name_for_user_and_dir}
Group=@@{name_for_user_and_dir}

Environment=GELDATA=@@{localstatedir}/lib/@@{name_for_user_and_dir}/@@{slot}/data/
RuntimeDirectory=@@{name_for_user_and_dir}

ExecStartPre=@@{pre_start_script}
ExecStart=@@{bindir}/@@{name} --data-dir=${GELDATA} --runstate-dir=%t/@@{name_for_user_and_dir} --tls-cert-mode=generate_self_signed
ExecReload=/bin/kill -HUP ${MAINPID}
KillMode=mixed
KillSignal=SIGINT
TimeoutSec=0

[Install]
WantedBy=multi-user.target
